# k8s-core/ansible/roles/bootstrap/ssh-known-hosts/files/renew-ssh-certificates.timer
[Unit]
Description=Renew SSH certificates from vault
After=network.target network-online.target

[Install]
WantedBy=multi-user.target

[Timer]
OnBootSec=10
# one day between renewals
OnUnitActiveSec=86400
